# homed.service
[Unit]
Description=homedaemon Service
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=homed.py -D

[Install]
WantedBy=multi-user.target

