# homed.service
[Unit]
Description=Smart home Service
Requires=couchdb.service
After=couchdb.service
After=network-online.target

[Service]
ExecStart=/usr/bin/homed.py

[Install]
Wants=network-online.target
