# homed.service
[Unit]
Description=Angry home Service
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/homed.py

[Install]
WantedBy=multi-user.target

