# housed.service
[Unit]
Description=Smart house Service

[Service]
ExecStart=/usr/bin/housed.py

[Install]
WantedBy=multi-user.target

